# dallas-student-navigator.service
[Unit]
Description=Dallas Student Navigator Streamlit App
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aws-agentic-ai-hackathon
Environment="PATH=/home/ubuntu/aws-agentic-ai-hackathon/venv/bin"
ExecStart=/home/ubuntu/aws-agentic-ai-hackathon/venv/bin/streamlit run app.py --server.address 0.0.0.0 --server.port 8501
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

